# Original Xiaomi Firmwares/Global-GB/3.0.24/miwifi_r4a_all_03233_3.0.24_INT/etc/hotplug.d/iface/81-speedtest
#!/bin/sh
[ "$ACTION" == "ifup" ] && [ "$INTERFACE" == "wan" ] && {
    /usr/sbin/auto_speedtestd oversea
}
