# Original Xiaomi Firmwares/Chinese-ZH-CN/2.28.62/miwifi_r4a_firmware_72d65_2.28.62/etc/init.d/smp
#!/bin/sh /etc/rc.common
# Copyright (C) 2006-2011 OpenWrt.org

START=19

start() {
	smp.sh wifi 2>&-
}
